# .config/systemd/user/pasystray.service
[Unit]
Description=Pulseaudio systray

[Service]
Type=simple
ExecStart=/bin/pasystray

[Install]
WantedBy=default.target
